# python-simple-server.service
[Unit]
Description=Python Simple Server

[Service]
User=jenkins
ExecStart=/opt/python-simple-server/app.py

[Install]
WantedBy=multi-user.target
